# Templates/LinuxPackageTemplate/debian/DEBIAN/postrm
#!/bin/sh

PATH=/sbin:/usr/sbin:/bin:/usr/bin

uninstall_daemon () {
  update-rc.d codesyscontrol remove
}

echo "uninstall daemon"
uninstall_daemon

exit 0
